# backend/supervisord.conf
[unix_http_server]
file=/home/vladis/project/sonanz/backend/var/run/supervisor.sock

[supervisord]
logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/supervisord.log
pidfile=/home/vladis/project/sonanz/backend/var/run/supervisord.pid
nodaemon=false

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///home/vladis/project/sonanz/backend/var/run/supervisor.sock

; --- Координация (1 воркер — намеренно, исключает гонку в allTracksHandled) ---
[program:sonanz-async]
command=/bin/php8.4 /home/vladis/project/sonanz/backend/bin/console messenger:consume async --time-limit=3600 --memory-limit=256M -vv
directory=/home/vladis/project/sonanz/backend
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=300
stdout_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/async.log
stderr_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/async.log

; --- Скачивание треков (параллельное yt-dlp) ---
[program:sonanz-downloads]
command=/bin/php8.4 /home/vladis/project/sonanz/backend/bin/console messenger:consume downloads --time-limit=3600 --memory-limit=512M -vv
directory=/home/vladis/project/sonanz/backend
numprocs=15
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=300
stdout_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/downloads.log
stderr_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/downloads.log

; --- Теггирование ---
[program:sonanz-tagging]
command=/bin/php8.4 /home/vladis/project/sonanz/backend/bin/console messenger:consume tagging --time-limit=3600 --memory-limit=256M -vv
directory=/home/vladis/project/sonanz/backend
numprocs=2
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=60
stdout_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/tagging.log
stderr_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/tagging.log

; --- Поиск ссылок ---
[program:sonanz-links]
command=/bin/php8.4 /home/vladis/project/sonanz/backend/bin/console messenger:consume links --time-limit=3600 --memory-limit=256M -vv
directory=/home/vladis/project/sonanz/backend
numprocs=1
process_name=%(program_name)s_%(process_num)02d
autostart=true
autorestart=true
startsecs=3
stopwaitsecs=60
stdout_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/links.log
stderr_logfile=/home/vladis/project/sonanz/backend/var/log/supervisor/links.log
